# 2016/march/systemd/getmail.service
[Unit]
Description=Getmail

[Service]
Type=simple
ExecStart=/sbin/getmail
